# etc/cron.d/crontab.conf
MAILTO=""
PHP=/usr/local/bin/php
WGET=/usr/bin/wget


*/5 * * * * $PHP /app/bin/console app:receive