# etc/cron.d/crontab.conf
*/5 * * * * $PHP /app/bin/console app:receive